# crontab
* * * * * php /home/Travis/www/prepevents_stage schedule:run >> /dev/null 2>&1
